# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-- \\u' --autologin user --noreset --noclear - $TERM
